# serversideup-php/8.3/fpm-nginx-alpine/etc/s6-overlay/s6-rc.d/laravel-queue/run
#!/command/execlineb -P
cd /var/www/html
php /var/www/html/artisan queue:work -v --tries=3
